# examples/mac-supervisor.ini
[program:piplin]
command=php artisan queue:work --queue=piplin-high,piplin-low --sleep=3 --tries=1 --daemon
directory=/data/wwwroot/piplin
process_name=%(program_name)s_queue_%(process_num)s
numprocs=1
numprocs_start=0
stdout_logfile=/var/log/supervisor/piplin-%(process_num)s-stdout.log
stderr_logfile=/var/log/supervisor/piplin-%(process_num)s-stderr.log
stderr_logfile_maxbytes=1MB
stdout_logfile_maxbytes=1MB
redirect_stderr=true
autostart=true
autorestart=true

[program:piplin-socketio]
command=/Users/hok/.nvm/versions/node/v8.16.0/bin/node socket.js
directory=/data/wwwroot/piplin
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
environment=HOME="/data/wwwroot/piplin",NODE_ENV="production"
stderr_logfile=/var/log/supervisor/piplin-socket-stderr.log
stdout_logfile=/var/log/supervisor/piplin-socket-stdout.log
stderr_logfile_maxbytes=1MB
stdout_logfile_maxbytes=1MB

[program:piplin-broadcast]
command=php artisan queue:work --queue=piplin-default --sleep=3 --tries=1 --daemon
directory=/data/wwwroot/piplin
numprocs=1
process_name=%(program_name)s_%(process_num)s
stderr_logfile=/var/log/supervisor/piplin-broadcast-stderr.log
stdout_logfile=/var/log/supervisor/piplin-broadcast-stdout.log
stderr_logfile_maxbytes=1MB
stdout_logfile_maxbytes=1MB
redirect_stderr=true
autostart=true
autorestart=true
